# deploy/systemd/echo-mcp-server.service
[Unit]
Description=Memory MCP Server
Documentation=https://github.com/Vladikasik/memory-mcp-server
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/memory-mcp-server
Environment=ENV=production
EnvironmentFile=/opt/memory-mcp-server/.env
ExecStart=/opt/memory-mcp-server/env/bin/python -m echo_mcp_server
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=5

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/memory-mcp-server
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target 